# wrangler/citydigest/wrangler/citydigest.service
[Unit]
Description=CityDigest News Service
After=network.target
StartLimitIntervalSec=0

[Service]
Type=simple
User=citydigest
WorkingDirectory=/home/citydigest/citydigest
ExecStart=/home/citydigest/venv/bin/python -m citydigest.main schedule
Environment="PYTHONUNBUFFERED=1"

# Email notification settings
Environment="SMTP_SERVER=smtp.gmail.com"
Environment="SMTP_PORT=587"
Environment="SMTP_USER=dakotacarrasco98@gmail.com"
Environment="SMTP_PASSWORD=odxr iqhk hwar lfjg"
Environment="NOTIFICATION_SENDER=dakotacarrasco98@gmail.com"
Environment="NOTIFICATION_RECIPIENT=dakota@newsr.io"

# Supabase settings
Environment="SUPABASE_URL=https://ckywgfnaaiapiyrztgdp.supabase.co"
Environment="SUPABASE_SERVICE_KEY=eyJhbGciOiJIUzI1NiIsInR5cCI6IkpXVCJ9.eyJpc3MiOiJzdXBhYmFzZSIsInJlZiI6ImNreXdnZm5hYWlhcGl5cnp0Z2RwIiwicm9sZSI6InNlcnZpY2Vfcm9sZSIsImlhdCI6MTc0Mzk3MDY0NCwiZXhwIjoyMDU5NTQ2NjQ0fQ.kR-jZNfDf8NfddReYx2irUVBh_3AYCLjeBM9UKvjGtw"
Environment="SUPABASE_KEY=eyJhbGciOiJIUzI1NiIsInR5cCI6IkpXVCJ9.eyJpc3MiOiJzdXBhYmFzZSIsInJlZiI6ImNreXdnZm5hYWlhcGl5cnp0Z2RwIiwicm9sZSI6ImFub24iLCJpYXQiOjE3NDM5NzA2NDQsImV4cCI6MjA1OTU0NjY0NH0.orM89MV8N4Nwdm9JmIWeVhu8-mqZqxdfp7VOZH-zM_s"
Environment="SUPABASE_ANON_KEY=eyJhbGciOiJIUzI1NiIsInR5cCI6IkpXVCJ9.eyJpc3MiOiJzdXBhYmFzZSIsInJlZiI6ImNreXdnZm5hYWlhcGl5cnp0Z2RwIiwicm9sZSI6ImFub24iLCJpYXQiOjE3NDM5NzA2NDQsImV4cCI6MjA1OTU0NjY0NH0.orM89MV8N4Nwdm9JmIWeVhu8-mqZqxdfp7VOZH-zM_s"
Environment="SUPABASE_JWT_SECRET=UCVV+q0GPvOQ6kjQ6OagmyyvCp8UD5QtRhr1blR8vL38BJoQhQUsYI5/1h7dIniS/xZQoF+wAj0j3QP2JNfByw=="
Environment="DEFAULT_AUTHOR_ID=00000000-0000-0000-0000-000000000000"

# Mixtral API settings
Environment="MIXTRAL_API_KEY=OTLFsDf7FEz44mcNREgxun2SDUHvGCVB"
Environment="MIXTRAL_API_ENDPOINT=https://api.mistral.ai/v1/chat/completions"

# Worker settings
Environment="WORKER_CONCURRENCY=3"
Environment="MAX_MEMORY_PERCENT=75"
Environment="LOG_LEVEL=INFO"
Environment="RUN_ON_STARTUP=true"

# Restart settings for crash recovery
Restart=always
RestartSec=10
RuntimeMaxSec=86400

# Output logs to journal
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target 